# conf/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
fsck.ext4 -y /dev/mmcblk0p2
exit 0
